# lighting-plot-container/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info

; Virtual display — must start before Perastage or x11vnc
[program:xvfb]
command=Xvfb :1 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
priority=10
autorestart=true
startsecs=1
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.log

; Kiosk window manager — must run before Perastage so it can claim the window
; -use_titlebar no  removes the WM title bar for a cleaner full-screen look
[program:wm]
command=matchbox-window-manager -use_titlebar no
environment=DISPLAY=":1"
priority=15
autorestart=true
startsecs=1
stdout_logfile=/var/log/supervisor/wm.log
stderr_logfile=/var/log/supervisor/wm.log

; Perastage application
; HOME points to the persistent volume so library/project data survives restarts.
; LIBGL_ALWAYS_SOFTWARE forces Mesa llvmpipe (no physical GPU needed).
; The wrapper checks /tmp/perastage-open-file and opens that path if present.
[program:perastage]
command=/opt/perastage-wrapper.sh
environment=
    DISPLAY=":1",
    HOME="/data/perastage/home",
    XDG_CONFIG_HOME="/data/perastage/home/.config",
    XDG_DATA_HOME="/data/perastage/home/.local/share",
    XDG_RUNTIME_DIR="/tmp/xdg-runtime",
    LIBGL_ALWAYS_SOFTWARE="1",
    MESA_GL_VERSION_OVERRIDE="4.5",
    MESA_GLSL_VERSION_OVERRIDE="450"
priority=20
autorestart=true
; Give Xvfb a moment to initialise before launching the GUI
startretries=5
startsecs=3
stdout_logfile=/var/log/supervisor/perastage.log
stderr_logfile=/var/log/supervisor/perastage.log

; VNC server attached to the virtual display
; -nopw: no password (rely on network/reverse-proxy auth instead)
; -forever: keep running after clients disconnect
; -shared: allow multiple simultaneous viewers
; -noxdamage: improves reliability with Xvfb (no DAMAGE extension)
[program:x11vnc]
command=x11vnc
    -display :1
    -nopw
    -listen localhost
    -rfbport 5900
    -forever
    -shared
    -xkb
    -cursor most
    -wait 50
priority=30
autorestart=true
startsecs=2
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc.log

; Webhook server — accepts "open file" commands from crew-hub on port 6081
[program:open-file-server]
command=python3 /opt/open-file-server.py
priority=25
autorestart=true
startsecs=1
stdout_logfile=/var/log/supervisor/open-file-server.log
stderr_logfile=/var/log/supervisor/open-file-server.log

; noVNC websocket bridge — exposes VNC over WebSocket on port 6080
; The noVNC static files (vnc.html etc.) are served from /usr/share/novnc/
[program:websockify]
command=websockify --web /usr/share/novnc 0.0.0.0:6080 localhost:5900
priority=40
autorestart=true
startsecs=2
stdout_logfile=/var/log/supervisor/websockify.log
stderr_logfile=/var/log/supervisor/websockify.log
